# home_users/root@hostcontainer/git/pwn1/main.socket
[Unit]
Description=tasty!
[Socket]
ListenStream=55
Accept=yes
[Install]
WantedBy=sockets.target
