# home_users/root@hostcontainer/git/pwn1/main@.service
[Unit]
Description=tasty!
[Service]
ExecStart=/bin/main
User=root
StandardInput=socket
